# sync.cron
* * * * * /opt/blocklist/sync_blocklist.sh >/dev/null 2>&1
